# cudavnc/vnc_bootstrap/supervisord.conf
[supervisord]
nodaemon = false
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid
directory = /tmp
childlogdir = /tmp

[program:Xvfb]
command = /usr/bin/Xvfb :0 -screen 0 1920x1080x16
autostart = true

[program:x11vnc]
command = /usr/bin/x11vnc -passwd atlas -display :0 -N -forever
autostart = true
